# launchd/com.jacob.granola-obsidian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jacob.granola-obsidian</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <!-- Update this path to wherever you cloned the repo -->
    <string>/path/to/granola_to_obsidian.py</string>
  </array>

  <!-- Runs daily at 9pm — adjust Hour/Minute as needed -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/Library/Logs/granola-obsidian.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/Library/Logs/granola-obsidian.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
